# launchd/LaunchDaemons/com.googlecode.munki.managedsoftwareupdate-install.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Debug</key>
  <true/>
  <key>Label</key>
  <string>com.googlecode.munki.managedsoftwareupdate-install</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/munki/managedsoftwareupdate</string>
    <string>--installwithnologout</string>
  </array>
  <key>WatchPaths</key>
  <array>
    <string>/Library/Managed Installs/.managedinstall.launchd</string>
  </array>
  <key>OnDemand</key>
  <true/>
</dict>
</plist>

